# share/jbt-desktop-fixer.service
[Unit]
Description=JetBrains Toolbox desktop fixer

[Service]
Type=simple
ExecStart={{exec}} watch

[Install]
WantedBy=default.target